# send-email/crontab
* * * * * /usr/local/bin/python /app/script.py
